# pypy/tool/build/bin/gentoo_init.d_webserver
#!/sbin/runscript

CMD=/opt/pypy-dist/pypy/tool/build/bin/webserver
PIDFILE=/var/run/build_webserver.pid
LOGFILE=/var/log/build_webserver.log

depend() {
	use net
}

start() {
	ebegin "Starting PyPy web server for build tool"
	start-stop-daemon --start --quiet --exec $CMD \
        --make-pidfile --pidfile $PIDFILE --background >> $LOGFILE
	eend $?
}

stop() {
	ebegin "Stopping PyPy web server for build tool"
	start-stop-daemon --stop --quiet --exec $CMD
	eend $?
}
